# docker_nginx/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:startnginx]
command=/start.sh

[program:fluentd]
command=/opt/utils/start-fluentd.sh
